# debian/postinst
#!/bin/bash
systemctl enable pardus-power-manager
